# iCloudMailFix/layout/Library/LaunchDaemons/com.victorlobe.iCloudMailFix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.victorlobe.iCloudMailFix</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/libexec/iCloudMailFixd</string>
  </array>

  <key>RunAtLoad</key><true/>

  <!-- Optional safety: slow down respawn if it exits instantly -->
  <key>ThrottleInterval</key><integer>3</integer>

  <!-- KeepAlive can be false; your PID lock already enforces single-instance -->
  <key>KeepAlive</key><false/>

  <key>UserName</key><string>root</string>
  <key>GroupName</key><string>wheel</string>

  <key>WorkingDirectory</key>
  <string>/usr/libexec</string>

  <key>StandardOutPath</key><string>/var/log/iCloudMailFixd.out</string>
  <key>StandardErrorPath</key><string>/var/log/iCloudMailFixd.err</string>
</dict>
</plist>